# deerflow-supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisor/supervisord.log

[program:langgraph]
command=bash -c "cd /app/backend && uv run langgraph dev --host 0.0.0.0 --port 2024 --no-browser --config /app/backend/langgraph.json"
autorestart=true
stdout_logfile=/var/log/langgraph.log
stderr_logfile=/var/log/langgraph_err.log
startsecs=3

[program:gateway]
command=bash -c "cd /app/backend && uv run uvicorn src.gateway.app:app --host 0.0.0.0 --port 8001"
autorestart=true
stdout_logfile=/var/log/gateway.log
stderr_logfile=/var/log/gateway_err.log
startretries=5
startsecs=5

[program:frontend]
command=bash -c "cd /app/frontend && npx next start -p 3000"
autorestart=true
stdout_logfile=/var/log/frontend.log
stderr_logfile=/var/log/frontend_err.log

[program:nginx]
command=nginx -g "daemon off;"
autorestart=true
stdout_logfile=/var/log/nginx_out.log
stderr_logfile=/var/log/nginx_err.log
